# overlay/etc/audit/rules.d/30-host-audit.rules
## "Log practically everything" — modeled on the Neo23x0 best-practice ruleset.
## Read back with:  ausearch -k <key> -i   |   aureport -f -i

## ===== Identity: user/group changes =====
-w /etc/passwd        -p wa -k identity
-w /etc/group         -p wa -k identity
-w /etc/shadow        -p wa -k identity
-w /etc/gshadow       -p wa -k identity
-w /etc/security/opasswd -p wa -k identity
-w /etc/nsswitch.conf -p wa -k identity

## ===== Sudo / privilege config =====
-w /etc/sudoers       -p wa -k actions
-w /etc/sudoers.d/    -p wa -k actions
-w /etc/polkit-1/     -p wa -k actions

## ===== Logins / logouts / sessions =====
-w /var/log/faillog   -p wa -k logins
-w /var/log/lastlog   -p wa -k logins
-w /var/run/utmp      -p wa -k session
-w /var/log/wtmp      -p wa -k session
-w /var/log/btmp      -p wa -k session
-w /etc/pam.d/        -p wa -k pam
-w /etc/ssh/sshd_config -p wa -k sshd

## ===== Command execution (what commands ran, by whom) =====
-a always,exit -F arch=b64 -S execve,execveat -F auid>=1000 -F auid!=unset -k process_creation
-a always,exit -F arch=b32 -S execve,execveat -F auid>=1000 -F auid!=unset -k process_creation

## ===== File modification: WHO edited WHAT and WHEN =====
## Permission / ownership / xattr changes
-a always,exit -F arch=b64 -S chmod,fchmod,fchmodat -F auid>=1000 -F auid!=unset -k perm_mod
-a always,exit -F arch=b32 -S chmod,fchmod,fchmodat -F auid>=1000 -F auid!=unset -k perm_mod
-a always,exit -F arch=b64 -S chown,fchown,fchownat,lchown -F auid>=1000 -F auid!=unset -k perm_mod
-a always,exit -F arch=b32 -S chown,fchown,fchownat,lchown -F auid>=1000 -F auid!=unset -k perm_mod
-a always,exit -F arch=b64 -S setxattr,lsetxattr,fsetxattr,removexattr,lremovexattr,fremovexattr -F auid>=1000 -F auid!=unset -k perm_mod
## Deletion / rename
-a always,exit -F arch=b64 -S unlink,unlinkat,rename,renameat,renameat2 -F auid>=1000 -F auid!=unset -k delete
-a always,exit -F arch=b32 -S unlink,unlinkat,rename,renameat,renameat2 -F auid>=1000 -F auid!=unset -k delete
## Failed file access under /etc (high signal)
-a always,exit -F arch=b64 -S open,openat,openat2,truncate,ftruncate -F dir=/etc -F success=0 -F auid>=1000 -F auid!=unset -k file_access_fail

## ===== Directory watches: every write under sensitive trees =====
-w /etc  -p wa -k etc_changes
-w /home -p wa -k home_changes

## ===== Scheduled tasks =====
-w /etc/crontab    -p wa -k cron
-w /etc/cron.d/    -p wa -k cron
-w /var/spool/cron/ -p wa -k cron
