# overlay/etc/systemd/system/opensnitch-prune.service
[Unit]
Description=Prune OpenSnitch event DB older than ~6 months

[Service]
Type=oneshot
ExecStart=/usr/local/bin/prune-opensnitch-db.sh
