# overlay/etc/systemd/system/opensnitch-prune.timer
[Unit]
Description=Daily OpenSnitch DB retention prune (~6 months)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
