# overlay/etc/systemd/system/zeek.service
[Unit]
Description=Zeek network security monitor (zeekctl standalone)
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
Environment=ZEEKPATH=/opt/zeek
# Resolve the primary (default-route) NIC at boot and inject it into node.cfg
# (interface name unknown at build time; done in a script so $ isn't eaten by systemd).
ExecStartPre=/usr/local/sbin/zeek-set-iface.sh
ExecStart=/opt/zeek/bin/zeekctl deploy
ExecStop=/opt/zeek/bin/zeekctl stop
RemainAfterExit=yes
Restart=on-failure

[Install]
WantedBy=multi-user.target
