# overlay/etc/systemd/system/zeek-cron.service
[Unit]
Description=zeekctl cron (log expiry + housekeeping)
After=zeek.service

[Service]
Type=oneshot
Environment=ZEEKPATH=/opt/zeek
ExecStart=/opt/zeek/bin/zeekctl cron
